# systemd/claude-fleet-orchestrate.timer
# Fires claude-fleet-orchestrate.service. Parity with the launchd plist:
# StartInterval=120 + RunAtLoad=true → run ~30s after the user manager comes up,
# then every 120s. Install: `systemctl --user enable --now
# claude-fleet-orchestrate.timer`.
[Unit]
Description=claude-fleet auto-orchestrator timer (every 120s)

[Timer]
OnStartupSec=30
OnUnitActiveSec=120
AccuracySec=10s

[Install]
WantedBy=timers.target
