# systemd/claude-fleet-collect.timer
# Fires claude-fleet-collect.service. Parity with the launchd plist:
# StartInterval=60 + RunAtLoad=true → run ~10s after the user manager comes up,
# then every 60s. Install: `systemctl --user enable --now
# claude-fleet-collect.timer`.
[Unit]
Description=claude-fleet collector timer (every 60s)

[Timer]
OnStartupSec=10
OnUnitActiveSec=60
AccuracySec=5s

[Install]
WantedBy=timers.target
