# systemd/claude-fleet-base-sync.timer
# Fires claude-fleet-base-sync.service. Parity with the launchd plist:
# StartInterval=60 + RunAtLoad=true.
[Unit]
Description=claude-fleet base-sync daemon timer (every 60s)

[Timer]
OnStartupSec=35
OnUnitActiveSec=60
AccuracySec=10s

[Install]
WantedBy=timers.target
