# systemd/claude-fleet-collect.service
# claude-fleet collector — Linux systemd user unit (parity with launchd
# com.claude-fleet.collect). Oneshot driven by claude-fleet-collect.timer.
[Unit]
Description=claude-fleet dashboard collector (git/gh/usage caches)

[Service]
Type=oneshot
ExecStart=/bin/bash __HOME__/.claude/fleet/bin/tmux-dash-collect.sh
Environment=HOME=__HOME__
Environment=LANG=en_US.UTF-8
Environment=PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
StandardOutput=null
StandardError=append:__HOME__/.claude/fleet/logs/collect.systemd.log
