# systemd/claude-fleet-cleanup.service
# claude-fleet cleanup daemon — Linux systemd user unit (parity with launchd
# com.claude-fleet.cleanup). Oneshot driven by claude-fleet-cleanup.timer.
# THE FLEET NEVER MERGES: /fleet-ship arms GitHub auto-merge; this daemon reaps
# the leftover worktree/window/branch and records the resume ledger once a PR is
# final (MERGED or CLOSED-unmerged). It merges nothing, so it is ON by default
# for every fleet (opt out per fleet with FLEET_CLEANUP=0). Needs gh.
[Unit]
Description=claude-fleet cleanup daemon (reap worktrees + record resume ledger after merges)

[Service]
Type=oneshot
ExecStart=/bin/bash __HOME__/.claude/fleet/bin/fleet-cleanup-daemon.sh
Environment=HOME=__HOME__
Environment=LANG=en_US.UTF-8
Environment=PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
StandardOutput=null
StandardError=append:__HOME__/.claude/fleet/logs/cleanup.systemd.log
