# systemd/claude-fleet-classify.timer
# Fires claude-fleet-classify.service. Parity with the launchd plist:
# StartInterval=1800 + RunAtLoad=true.
# Backstop only: the Stop hook (classify-hook.sh) does real-time single-window
# classification; this full scan just catches windows no Stop will revisit.
[Unit]
Description=claude-fleet classifier backstop timer (every 1800s)

[Timer]
OnStartupSec=30
OnUnitActiveSec=1800
AccuracySec=60s

[Install]
WantedBy=timers.target
