# systemd/claude-fleet-summarize.service
# claude-fleet summarizer — Linux systemd user unit (parity with launchd
# com.claude-fleet.summarize). Oneshot driven by claude-fleet-summarize.timer.
# Optional: spends small, change-gated LLM tokens. Needs `claude` on PATH.
[Unit]
Description=claude-fleet session summarizer (1-line LLM summary per session)

[Service]
Type=oneshot
ExecStart=/bin/bash __HOME__/.claude/fleet/bin/tmux-summarize.sh
Environment=HOME=__HOME__
Environment=LANG=en_US.UTF-8
Environment=PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
StandardOutput=null
StandardError=append:__HOME__/.claude/fleet/logs/summarize.systemd.log
