# systemd/claude-fleet-dispatch.service
# claude-fleet autofill dispatcher — Linux systemd user unit (parity with launchd
# com.claude-fleet.dispatch). Oneshot driven by claude-fleet-dispatch.timer.
# OFF by default: only fleets whose conf sets FLEET_AUTOFILL=1 auto-spawn, and
# auto-spawning launches real Claude sessions that spend LLM tokens. Needs gh.
[Unit]
Description=claude-fleet autofill dispatcher (auto-spawn backlog by priority under caps)

[Service]
Type=oneshot
ExecStart=/bin/bash __HOME__/.claude/fleet/bin/fleet-dispatch.sh
Environment=HOME=__HOME__
Environment=LANG=en_US.UTF-8
Environment=PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
StandardOutput=null
StandardError=append:__HOME__/.claude/fleet/logs/dispatch.systemd.log
