# launchd/com.claude-fleet.webhook.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.claude-fleet.webhook</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string><string>__HOME__/.claude/fleet/bin/fleet-webhook.sh</string>
  </array>
  <key>EnvironmentVariables</key><dict>
    <key>HOME</key><string>__HOME__</string>
    <key>LANG</key><string>en_US.UTF-8</string>
    <key>PATH</key><string>__HOME__/.local/bin:__BREW_PREFIX__/bin:__BREW_PREFIX__/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ProcessType</key><string>Background</string>
  <key>StandardErrorPath</key><string>__HOME__/.claude/fleet/logs/webhook.launchd.log</string>
  <key>StandardOutPath</key><string>/dev/null</string>
</dict></plist>
