# launchd/com.claude-fleet.spinner.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.claude-fleet.spinner</string>
  <key>ProgramArguments</key><array>
    <string>/bin/sh</string><string>__HOME__/.claude/fleet/bin/tmux-spinner.sh</string>
  </array>
  <key>EnvironmentVariables</key><dict>
    <key>PATH</key><string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>LANG</key><string>en_US.UTF-8</string>
    <key>SPIN_INTERVAL</key><string>0.12</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ProcessType</key><string>Background</string>
  <key>StandardErrorPath</key><string>__HOME__/.claude/fleet/logs/spinner.launchd.log</string>
  <key>StandardOutPath</key><string>/dev/null</string>
</dict></plist>
